# gammu-disc.service
[Unit]
Description=gammu-disc
After=gammu-smsd.service

[Service]
ExecStart=/opt/gammu-disc/gammu-disc

[Install]
WantedBy=multi-user.target